# variants/multi-agent/install/host-side/gf-notify-architect.service
[Unit]
Description=GateForge — dispatch signed notification to Architect
Documentation=https://github.com/tonylnng/gateforge-openclaw-configs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/opt/secrets/gateforge.env
ExecStart=/opt/gateforge/bin/gf-notify-architect.sh
Nice=5
TimeoutStartSec=60

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/gateforge /tmp
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
